# SRC/DEB_picopilot-idrop/etc/supervisor/conf.d/kBlue.conf
[program:kBlue]
autostart=false
command=/usr/bin/python3 /opt/piCopilot-idrop/kBlue.py --pgsql
directory=/opt/piCopilot-idrop
stderr_logfile=/tmp/kBlueERR
stderr_logfile=NONE
