# SRC/DEB_picopilot-idrop/etc/supervisor/conf.d/kBlue.conf
[program:kBlue]
autostart=false
command=/usr/bin/ubertooth-btle -f -q /mnt/usb_storage/bluesPipe 1>/dev/null
directory=/opt/piCopilot-idrop
stderr_logfile=/tmp/kBlueERR
stdout_logfile=/tmp/kBlueOUT
